# app/conf.d/websockify.conf
[program:websockify]
command=websockify --web /novnc %(ENV_NOVNC_WEB_PORT)s localhost:5901
autorestart=true
